# dot_config/systemd/user/qmd-embed.timer
[Unit]
Description=QMD vector embedding rebuild daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
