# dot_config/systemd/user/qmd-embed.service
[Unit]
Description=QMD vector embedding rebuild (slow)

[Service]
Type=oneshot
ExecStart=/home/hschne/.local/share/mise/shims/qmd embed
